# bash/devops/templates/supervisor/lk/exchange.ini
[program:exchange_#HOST_NAME#]
command = bash -c "ulimit -n 10000; exec php /var/www/#HOST_NAME#/api/yii exchange/actions/start"
process_name = exchange_#HOST_NAME#
numprocs = 1
autostart = true
autorestart = true
user = nginx
stdout_logfile = /var/log/supervisor/exchange_#HOST_NAME#.info.log
stdout_logfile_maxbytes = 10MB
stderr_logfile = /var/log/supervisor/exchange_#HOST_NAME#.error.log
stderr_logfile_maxbytes = 100MB
stopsignal=TERM
stopwaitsecs=3
stopasgroup=true
killasgroup=true
