# bash/devops/templates/supervisor/lk/notifying.ini
[program:notifying_#HOST_NAME#]
command = bash -c "ulimit -n 10000; exec php /var/www/#HOST_NAME#/api/yii notification/run"
process_name = notifying_#HOST_NAME#
numprocs = 1
autostart = true
autorestart = true
user = root
stdout_logfile = /var/log/supervisor/notifying_#HOST_NAME#.info.log
stdout_logfile_maxbytes = 5MB
stderr_logfile = /var/log/supervisor/notifying_#HOST_NAME#.error.log
stderr_logfile_maxbytes = 100MB
stopsignal=TERM
stopwaitsecs=5
stopasgroup=true
killasgroup=true